# playbooks/packages/build/templates/rootfs/usr/lib/systemd/system/sshd.service.j2
[Unit]
Description=OpenSSH Daemon
After=network.target

[Service]
ExecStartPre=/usr/bin/ssh-keygen -A
ExecStart=/usr/sbin/sshd -D
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
